# deploy/systemd/emaillm.service.d/env.conf
[Service]
EnvironmentFile=/srv/emaillm-ws/.env
Environment="PYTHONPATH=/srv/emaillm-ws/src"